# pgweb/pgweb.conf
[program:pgweb]
autostart=true
autorestart=true
command=/usr/bin/pgweb --bind %(ENV_PGWEB_BIND)s --port %(ENV_PGWEB_PORT)s --sessions
